# remmina-time-tracker.service
[Unit]
Description=Remmina Time Tracker - Session logging and idle auto-disconnect
Documentation=https://github.com/Agneskoinange/remmina-time-tracker
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/remmina-time-tracker
Restart=on-failure
RestartSec=5
# DISPLAY is auto-detected by the launcher script
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remmina-time-tracker

[Install]
WantedBy=graphical-session.target
